# ktox-with-m5.service
[Unit]
Description=KTOX_Pi with M5Cardputer Remote Control Support
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/KTOX_Pi
Environment="RJ_FRAME_MIRROR=1"
Environment="RJ_FRAME_PATH=/dev/shm/ktox_last.jpg"
Environment="RJ_FRAME_FPS=6"
Environment="RJ_WS_HOST=0.0.0.0"
Environment="RJ_WS_PORT=8765"
ExecStart=/usr/bin/python3 ktox_device_root.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
